# net-fs:nfs-utils/proc-fs-nfsd.mount
[Unit]
Description=nfsd file system
DefaultDependencies=no

[Mount]
What=nfsd
Where=/proc/fs/nfsd
Type=nfsd
