# net-fs:nfs-utils/rpc.gssd.service
[Unit]
Description=NFS rpcsec_gss daemon
Requires=var-lib-nfs-rpc_pipefs.mount
Requires=rpcbind.target
After=syslog.target var-lib-nfs-rpc_pipefs.mount rpcbind.target
Before=remote-fs.target

[Service]
ExecStart=/usr/sbin/rpc.gssd -f

[Install]
WantedBy=multi-user.target
